# gtservice.service
[Unit]
Description=Example GeoTrellis Server

[Service]
Type=forking
ExecStart=/bin/bash /usr/bin/gtservice.sh

[Install]
WantedBy=multi-user.target
